# bering-uclibc4/source/vsftpd/vsftpd.cron
#!/bin/sh
# Save daily LOGDEPTH versions of syslogfile
LOGDEPTH=4
# syslogfile name
LOGFILE=/var/log/vsftpd.log

if [ -f $LOGFILE ]; then
	savelog -g wheel -m 644 -u root -c $LOGDEPTH $LOGFILE >/dev/null
	/etc/init.d/vsftpd reload
fi
